# tools/dnc/systemd/dnc.service.in
[Unit]
Description=Dnc server
Documentation=man:dnc(1)
After=network.target

[Service]
Type=exec
User=dnc
Group=dnc
ExecStart=@bindir@/dnc -l -d @sysconfdir@/dhtnet/dnc.yaml -c @sysconfdir@/dhtnet/id/id-server.crt -p @sysconfdir@/dhtnet/id/id-server.pem
Restart=on-failure
ProtectSystem=strict
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
PrivateDevices=yes
PrivateUsers=yes
SystemCallFilter=@system-service

[Install]
WantedBy=multi-user.target
